# chirp2mqtt/rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/usr/bin/env bash
#set -x
#echo "++ rest api ++"

#chirpstack-rest-api
exec chirpstack-rest-api --insecure --cors "*"

